# flask_blog.service
[Unit]

[Service]
ExecStart=/root/flask_blog/app.py

[Install]
